# batnag.service
[Unit]
Description=Low battery indicator

[Service]
ExecStart=/usr/local/bin/batnag
Restart=always

[Install]
WantedBy=multi-user.target
