# compton.service
[Unit]
Description=Compositor for X11
After=graphical.target

[Service]
ExecStart=/usr/bin/compton

[Install]
WantedBy=graphical.target
