# xmodmap.service
[Unit]
Description=Customized keymap
After=graphical.target

[Service]
ExecStart=/usr/bin/xmodmap %h/.Xmodmap

[Install]
WantedBy=graphical.target
